# config/supervisord.conf
[program:graffitiwall]
command=/usr/bin/node server.js
directory=/var/www/graffitiwall/graffitiwall/server
autostart=true
autorestart=true
startsecs=10
startretries=3
user=yarek
