# config/walld.conf
[program:walld]
command=/var/www/graffitiwall/server/server.php
directory=/var/www/graffitiwall/server
autostart=true
autorestart=true
startsecs=10
startretries=3
user=yarek
